# systemd/backup.service
[Unit]
Description=rsync a few important folders to 320GB drive

[Service]
Type=oneshot
ExecStart=/home/michael/MEGA/bin/backup.sh

