# systemd/backup.timer
[Unit]
Description=Execute backup ce per day

[Timer]
OnCalendar=*-*-* 08,18:00:00
Unit=backup.service

[Install]
WantedBy=timers.target
