# systemd/configScripts.service
[Unit]
Description=rsync configs and scripts

[Service]
Type=oneshot
ExecStart=/home/michael/MEGA/bin/configScripts.sh

