# systemd/rclone.service
[Unit]
Description=rclone stuff to Dropbox site

[Service]
Type=oneshot
ExecStart=/home/michael/MEGA/bin/rclone.sh

