# systemd/shtf.service
[Unit]
Description=rsync stuff to shtf folder

[Service]
Type=oneshot
ExecStart=/home/michael/MEGA/bin/shtf.sh

