# systemd/termuxTo.service
[Unit]
Description=Sync to termux

[Service]
Type=oneshot
ExecStart=/home/michael/MEGA/bin/termuxTo.sh

