# systemd/termuxFrom.service
[Unit]
Description=Sync from termux

[Service]
Type=oneshot
ExecStart=/home/michael/MEGA/bin/termuxFrom.sh

