# systemd/termuxTo.timer
[Unit]
Description=Execute backup every day every hour

[Timer]
OnCalendar=*-*-* 00/1:35:00
Unit=termuxTo.service

[Install]
WantedBy=timers.target
